# systemd-timer/ssl-renew.timer
[Unit]
Description=Run ssl-renew.service daily at midnight

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target